# overlay/etc/periodic/hourly/php-clean-sessions
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Call php-clean-sessions executable
#======================================================================================================================

s6-setuidgid www php-clean-sessions
