# overlay/etc/periodic/hourly/php-clean-sessions
#!/usr/bin/nu

use bf-nginx-php sessions clean ; clean
