# overlay/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-nginx-php run ; run preflight"


#======================================================================================================================
# Start PHP FastCGI Process Manager.
#
# Flags:
#   -F  run in foreground mode
#======================================================================================================================

php-fpm -F
