# overlay/etc/services.d/php-error-log/run
#!/usr/bin/with-contenv bash
set -euo pipefail
_forward "php-error-log" "${PHP_INI_ERROR_LOG}"
