# overlay/etc/services.d/php-error-log/run
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Forward PHP errors to Docker
#======================================================================================================================

bcg-forward "php-error-log" "${PHP_INI_ERROR_LOG}"
