# overlay/etc/services.d/php-error-log/run
#!/usr/bin/with-contenv bash
set -euo pipefail

[[ -z "${PHP_INI_ERROR_LOG}" ]] && _disable "php-error-log"

if [ -e ${PHP_INI_ERROR_LOG} ]; then
    _forward "${PHP_INI_ERROR_LOG}"
else
    # wait for service to start and create log file
    exec sleep 5
fi