# overlay/etc/services.d/php-fpm/finish
#!/usr/bin/with-contenv bash

set -euo pipefail


#======================================================================================================================
# Use base executable to terminate all services
#======================================================================================================================

_terminate
