# overlay/etc/services.d/php-fpm/finish
#!/usr/bin/with-contenv bash
_terminate
