# overlay/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bash
_echo "starting PHP in FPM mode"
php-fpm7 -F
