# overlay/etc/services.d/php-fpm/run
#!/bin/bash

echo " - starting PHP-FPM"
php-fpm7 -F
