# overlay/etc/services.d/php-fpm/run
#!/bin/sh

echo " - starting PHP in FPM mode"
php-fpm7 -F
